# hosts/armhost-cm4/server/si_bmc.service
[Unit]
Description=SI BMC Server (Orange Pi CM4 KVM/BMC)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/si-bmc
# The "-" prefix means: ignore if this command fails
# (USB Gadget may not be available if UDC is not configured)
ExecStartPre=-/opt/si-bmc/setup_gadget.sh setup
ExecStart=/usr/bin/python3 /opt/si-bmc/main.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Security
NoNewPrivileges=false
ProtectSystem=false
ReadWritePaths=/sys/class/gpio /sys/kernel/config /dev

[Install]
WantedBy=multi-user.target
